# apps/cli/systemd/grota-verify.service
[Unit]
Description=Grota weekly backup verification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=grota
Group=grota
EnvironmentFile=/etc/grota/grota.env
ExecStart=/usr/local/bin/grota each audit backup
TimeoutStartSec=7200
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/backup /var/log/grota /var/lock/grota /tmp
PrivateTmp=true
