# apps/cli/systemd/grota-backup.service
[Unit]
Description=Grota backup - all accounts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=grota
Group=grota
EnvironmentFile=/etc/grota/grota.env
ExecStart=/usr/local/bin/grota backup all
TimeoutStartSec=14400
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/backup /var/log/grota /var/lock/grota /tmp
PrivateTmp=true
